# smh-api/smh_service/etc/smh_api.ini
[supervisord]
environment =
  AIS_REST_BASE_URL="http://internal-prod-ais-elasticlo-1et06hcch616u-2099117978.us-east-1.elb.amazonaws.com/api/v2",
  AIS_USERNAME="admin",
  AIS_PASSWORD="secret",
  DB_HOST="stage-rds-smh.cruttu7ke9i4.us-east-1.rds.amazonaws.com",
  DB_NAME="smh_service",
  DB_USER="smh_api",
  DB_PASSWORD="noan7wiu9OoF",
  LOG_LEVEL="INFO",
  ENVIRONMENT="test",
  PORT_SERVICE_BASE_URL="ports.polestar-testing.local:50051",
  PYTHONPATH = "/opt/smh",
  SIS_BASE_URL="https://sis.polestar-testing.com/api/v1",
  SIS_USERNAME="purpletrac",
  SIS_API_KEY="b85d85a9523e7caf3ecec776b003f3c91e84aa94",
  STATSD_HOST="",
  STATSD_PREFIX="smh"

[program:smh]
user = devops
directory = /opt/smh/smh_service
command = uwsgi --http :8080 --wsgi-file smh_api.py --callable app --processes 8 --http-timeout 600
autostart = true
autorestart = true
stdout_logfile = /var/log/polestar/smh/smh_stdout.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups=10
stderr_logfile = /var/log/polestar/smh/smh_stderr.log
stopasgroup = true
killasgroup = true
process_name = %(program_name)s
